# actions/com.jarvis.end-of-day.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Run end-of-day (wrap-up + "ok, Im leaving" prompt) every workday at 6 PM.
  Install: cp actions/com.jarvis.end-of-day.plist ~/Library/LaunchAgents/ && launchctl load ~/Library/LaunchAgents/com.jarvis.end-of-day.plist
  Uninstall: launchctl unload ~/Library/LaunchAgents/com.jarvis.end-of-day.plist
  Update the script path in ProgramArguments if the Jarvis repo is moved.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jarvis.end-of-day</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/osascript</string>
    <string>-e</string>
    <string>tell application &quot;Terminal&quot; to do script &quot;/Users/ricardosolis/GitHub/Jarvis/actions/end-of-day.sh&quot;</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ricardosolis/GitHub/Jarvis</string>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/com.jarvis.end-of-day.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.jarvis.end-of-day.stderr.log</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
  </array>
</dict>
</plist>
